# vault-observer.service
# =============================================================================
# vault-observer.service — systemd user service unit
#
# Install location: ~/.config/systemd/user/vault-observer.service
#
# Enable & start:
#   systemctl --user daemon-reload
#   systemctl --user enable --now vault-observer.service
#
# View logs:
#   journalctl --user -u vault-observer.service -f
#   # or the script's own log:
#   tail -f ~/.local/logs/vault-observer.log
# =============================================================================

[Unit]
Description=Vault Observer — Obsidian Git Auto-Commit Daemon
Documentation=https://github.com/your-username/vault-observer
After=network.target graphical-session.target
# Restart if it crashes — will not restart on clean exit (code 0)
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
ExecStart=%h/.local/bin/vault-observer.sh
Restart=on-failure
RestartSec=15s

# ── Environment ─────────────────────────────────────────────────────────────
# Edit these to match your setup.
Environment=VAULT_DIR=%h/obsidian-vault
Environment=INTERVAL_SECONDS=180
Environment=LOG_FILE=%h/.local/logs/vault-observer.log
Environment=GIT_BRANCH=main
Environment=COMMIT_AUTHOR_NAME=Vault Observer
Environment=COMMIT_AUTHOR_EMAIL=observer@local

# ── Resource limits (lightweight by design) ─────────────────────────────────
# CPUQuota: cap at 5% of one core — inotifywait is event-driven, not polling
CPUQuota=5%
# MemoryMax: git + inotify + zsh use ~20MB at peak; 64MB is very generous
MemoryMax=64M
# Nice: run at low priority so it never competes with Obsidian itself
Nice=15
IOSchedulingClass=idle
IOSchedulingPriority=7

# ── Hardening (optional but recommended) ────────────────────────────────────
# Remove if you encounter issues; these restrict what the process can do.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=%h/obsidian-vault %h/.local

[Install]
WantedBy=default.target
